# websh.service
[Unit]
Description=websh — SSH terminal backend
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/websh/server.py
Restart=on-failure
RestartSec=5

Environment=PORT=8765
Environment=HOST=127.0.0.1
Environment=SESSION_TIMEOUT=300
Environment=MAX_SESSIONS=10

User=websh
Group=websh

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ReadWritePaths=/tmp

# Open-file limit. Each active session holds: 1 PTY master_fd + 1 SSH
# ControlMaster socket + 1 epoll fd in every SSE/long-poll worker. Under
# 50 sessions × split panes + side-channel uploads the default 1024 is
# easy to exhaust without a clear error. 4096 leaves plenty of slack.
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
