# websh@.service
# Template unit: one websh instance per /etc/websh/<name>.env.
#
#   install -d /etc/websh
#   printf 'PORT=8766\n' > /etc/websh/staging.env
#   systemctl enable --now websh@staging
#
# Instances are fully isolated: own port, own state directory
# (/var/lib/websh-<name>), own env file. The single-instance
# websh.service remains the default path for one-box deployments.
#
# NOTE: derived from websh.service — apply hardening/option changes
# to BOTH files.
[Unit]
Description=websh — SSH terminal backend (instance %i)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/websh/server.py
Restart=on-failure
RestartSec=5

# Per-instance configuration lives in /etc/websh/<instance>.env —
# PORT is the one knob every instance MUST set (instances collide on
# the default otherwise). Everything else (HOST, SESSION_TIMEOUT,
# MAX_SESSIONS, WEBSH_CONFIG, WEBSH_VAULT_ENABLE, ...) is optional;
# see docs/configuration.md for the full table.
Environment=HOST=127.0.0.1
Environment=SESSION_TIMEOUT=300
Environment=MAX_SESSIONS=50
EnvironmentFile=/etc/websh/%i.env

# Per-instance writable state (the credential vault lives here when
# enabled: WEBSH_CREDS_PATH=/var/lib/websh-%i/websh.creds.json).
StateDirectory=websh-%i

# The encrypted credential vault is available but OFF by default. To turn
# it on, `systemctl edit websh` and add `Environment=WEBSH_VAULT_ENABLE=1`
# (the cryptography wheel from requirements.txt must be installed; without
# it the server still runs and the saved-credential UI just stays hidden).
# StateDirectory pre-creates /var/lib/websh owned by the service user and
# keeps it writable despite ProtectSystem=strict / ProtectHome=read-only,
# so the lazily created websh.creds.json (mode 0600) has somewhere to live
# the moment the vault is enabled.
Environment=WEBSH_CREDS_PATH=/var/lib/websh/websh.creds.json
StateDirectory=websh

User=websh
Group=websh

# Security hardening. websh drives ssh from web-supplied input, so the
# blast radius of a backend compromise is worth shrinking aggressively.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes

# The service needs no privileges of its own — it only spawns ssh as the
# unprivileged service user.
CapabilityBoundingSet=
AmbientCapabilities=

# ssh needs outbound IPv4/IPv6; AF_UNIX is required for the ControlMaster
# socket websh keeps under (private) /tmp. AF_NETLINK is required too:
# glibc's getaddrinfo() opens a NETLINK_ROUTE socket (__check_pf) to
# enumerate local address families, so without it DNS resolution can fail
# — both in server.py's pre-connect _resolve_host_ips() (whose gaierror
# fall-open would then silently disable the denied-hosts SSRF check) and
# in the ssh children. Everything else (AF_PACKET, AF_BLUETOOTH, …) is
# denied.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK

# Allow only the syscall groups a normal service (and its ssh children)
# use. @system-service already permits fork/exec/wait, so spawning ssh
# still works; raw/exotic syscalls are blocked.
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native

# CPython has no JIT and ssh maps no W+X pages, so write-xor-execute is safe.
MemoryDenyWriteExecute=yes

# Kernel / host surface the service has no business touching.
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=pid
PrivateDevices=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
UMask=0077

# Open-file limit. Each active session holds: 1 PTY master_fd + 1 SSH
# ControlMaster socket + 1 epoll fd in every SSE/long-poll worker. Under
# 50 sessions × split panes + side-channel uploads the default 1024 is
# easy to exhaust without a clear error. 4096 leaves plenty of slack.
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
